# dist/init/systemd/bumblebee-configurator.service
[Unit]
Description=Bumblebee configurator
After=cloud-init.service

[Service]
Type=oneshot
ExecStart=/opt/bumblebee/scripts/configurator
RemainAfterExit=yes
TimeoutSec=0
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
